# weathers.timer
[Unit]
Description=Timer to automatically run weathers service

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
